# hugo/hugo-build.container
[Container]
Image=docker.io/klakegg/hugo:ext-alpine
ContainerName=hugo-build
Volume=/path/to/website:/src:z

[Service]
Restart=no
