# focalboard/focalboard.container
[Container]
Image=docker.io/mattermost/focalboard
ContainerName=focalboard
AutoUpdate=registry

PublishPort=8000:8000

[Service]
Restart=on-success
