# wekan/wekan-db.container
[Container]
Image=docker.io/mongo
ContainerName=wekan-db

Network=wekannet.network

[Service]
Restart=on-success
