# wekan/wekan-db.container
[Container]
Image=docker.io/mongo
ContainerName=wekan-db
AutoUpdate=registry
Notify=true

Network=wekannet.network

[Service]
Restart=on-success
